# zookeeper_supervisor.conf
[program:zookeeper]
command=/opt/kafka_2.11-0.10.0.0/bin/zookeeper-server-start.sh /opt/kafka_2.11-0.10.0.0/config/zookeeper.properties
autostart=true
autorestart=true